# on_boot.service
[Unit]
Description=start-fan on boot.

[Service]
Type=oneshot
ExecStart=/usr/local/bin/start_fan.sh

[Install]
WantedBy=multi-user.target

