# on_boot.service
[Unit]
Description=start-fan on boot.

[Service]
Type=oneshot
ExecStart=/git/PI5/systemd/start_fan.sh

[Install]
WantedBy=multi-user.target

