# systemd/git_pull.service
[Unit]
Description=Run git_pull scrip

[Service]
Type=oneshot
ExecStart=/usr/local/bin/git_pull.sh
User=root
